# kurento-tree-demo/package/support-files/systemd.service
[Unit]
Description=${project.description}

[Service]
User=nobody
ExecStart=/var/lib/kurento/${project.artifactId}.jar run

[Install]
WantedBy=multi-user.target
